# scripts/com.kaplansolutions.outreach-awake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kaplansolutions.outreach-awake</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__ROOT__/scripts/outreach_keep_awake.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__ROOT__</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__ROOT__/data/outreach-awake.log</string>
    <key>StandardErrorPath</key>
    <string>__ROOT__/data/outreach-awake.log</string>
</dict>
</plist>
